# git-check.cron
# git-check.cron — Daily git repository health check
#
# Drop this file in /etc/cron.d/ to activate:
#   sudo cp ~/dev-env/git-check.cron /etc/cron.d/git-daily-check
#
# Runs git-check.sh inside the claude-dev container via docker exec.
# Output appends to ~/git-check.log on the host.
#
# Format: minute hour day-of-month month day-of-week user command
0 8 * * * ubuntu docker exec claude-dev bash /home/dev/dev-env/git-check.sh >> /home/ubuntu/git-check.log 2>&1
